# bp/supervisord.conf
[program:node]
command = node /home/dotcloud/current/front.js
directory = /home/dotcloud/current/